# test/Docker/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/sbin/mysqld --datadir=/var/lib/mysql --user=mysql
autostart=true
autorestart=true
priority=1

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
autorestart=true
priority=2
